# nix/systemd/nft-reload.service
[Unit]
Description=Reload nftables ruleset

[Service]
Type=oneshot
ExecStart=@nft@/sbin/nft -f /etc/nftables.conf